# plugins/slack-bridge/autostart-templates/linux-systemd.service
[Unit]
Description=Claude Slack Bridge daemon
After=network.target

[Service]
Type=simple
ExecStart=${NODE_PATH} ${BRIDGE_ENTRY}
WorkingDirectory=${CONFIG_DIR}
Restart=on-failure
RestartSec=10
StandardOutput=append:${LOG_DIR}/bridge.stdout.log
StandardError=append:${LOG_DIR}/bridge.stderr.log

[Install]
WantedBy=default.target
